# .config/supervisor/uwsgi.conf
[program:uwsgi]
command=/root/.pyenv/versions/onoff_crawling/bin/uwsgi --ini /srv/onoff_crawling/.config/uwsgi/uwsgi-app.ini

